# curator/crontabs/root
* * * * * curator /delete_index.yml --config /curator.yml --dry-run 
